# openwrt_pybuilder/resources/templates/base/files/etc/uci-defaults/99_custom-passwd.sh
#!/bin/ash
#
# Copyright (C) 2019 xxy1991
#

echo 'root:OW-Password' | chpasswd -m

exit 0
